# pkg/osx/scripts/postinstall
#!/bin/sh
set -e
/bin/launchctl load "/Library/LaunchDaemons/com.qosient.argus.plist"
